# DEBIAN/postinst
#! /bin/sh
# postinst script for wlanthermo

# enable all required services
systemctl enable WLANThermo.service
systemctl enable WLANThermoWD.service
systemctl enable pigpiod.service
systemctl enable nextion_hwclock.service

# restart sudo service
systemctl restart sudo 

# enable lighttpd module
lighty-enable-mod auth

if [ -f /etc/lighttpd/lighttpd.conf ]
then
    echo Backing up old lighttpd.conf
    mv /etc/lighttpd/lighttpd.conf /etc/lighttpd/lighttpd.conf.bak
fi

mv /var/www/conf/lighttpd.conf /etc/lighttpd/lighttpd.conf

# check if new web password has been set
if [ -f /var/www/tmp/update ]
then
    echo 'Password found!'
    passwort=`cat /var/www/tmp/update`
    echo "$passwort"
    htpasswd -c -b /etc/lighttpd/htpasswd wlanthermo $passwort
    rm /var/www/tmp/update # Updateflag zuruecksetzen
else
    if [ -t 0 ]
    then
        echo 'No password found!'
        echo 'Please enter the new password for user >wlanthermo< of the protected area!'
        htpasswd -c /etc/lighttpd/htpasswd wlanthermo
    else
        echo 'No password found!'
        echo 'Setting default password'
        htpasswd -c -b /etc/lighttpd/htpasswd wlanthermo raspberry
    fi
fi

ln -s /var/log/WLAN_Thermo /var/www/thermolog
ln -s /var/log/WLANThermo.log /var/www/log/WLANThermo.log
chown -Rv www-data:www-data /var/www 1>/dev/null
chown -Rv www-data:www-data /var/log/WLAN_Thermo 1>/dev/null

chown -v root:root /etc/sudoers.d/WLANThermo 1>/dev/null
chmod 0440 /etc/sudoers.d/WLANThermo 1>/dev/null

# Reload new unit files
systemctl daemon-reload

# Restart services
systemctl restart lighttpd.service
systemctl restart WLANThermo.service
systemctl start pigpiod.service
# systemctl restart WLANThermoWD.service
