# software/lib/systemd/system/WLANThermoUPDATECHECK.timer
[Unit]
Description=Run WLANThermo update check daily and after boot

[Timer]
OnBootSec=5min
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target