# software/etc/cron.daily/wlanthermo_updatecheck
#!/bin/sh
/usr/bin/update_checker.py wlanthermo > /var/www/tmp/updates.json